# my_files/etc-files/uci-defaults/95-mtk-led-fix-enable
#!/bin/sh
# Enable and run the BPI-R4 port-LED register fix at boot.
[ -x /etc/init.d/mtk-led-fix ] || exit 0
/etc/init.d/mtk-led-fix enable
/etc/init.d/mtk-led-fix start
exit 0
